# scripts/unison/unison-csync.timer
[Unit]
Description=Run Unison C Sync every 3 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target